# tgbot.service.example
[Unit]
Description=AI Server Admin Telegram Bot
After=network.target

[Service]
Type=simple
WorkingDirectory=/root/tgbot
ExecStart=/usr/bin/python3 bot.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
